# cron/crontab
* * * * * /app/rake.sh tools:helloworld
